# debian-live/source/config/includes.chroot_after_packages/etc/systemd/system/serial-getty@ttyS1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --keep-baud --autologin installer 115200 - $TERM
